# daemon/window-watcher.service
[Unit]
Description=LifeLogger Window Watcher
After=lifelogger.service

[Service]
Type=simple
ExecStart=%h/.local/share/lifelogger/venv/bin/python window_watcher.py
WorkingDirectory=%h/.config/lifelogger/daemon
Restart=always
RestartSec=5
Environment=DISPLAY=:0
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
